# packaging/debian/postrm
#!/bin/sh
# postrm: post-remove script for tui-serve
# Follows dh_installsystemd canonical patterns (Debian Policy §9.3)
set -e

DEB_SYSTEMD_HELPER=""
if command -v deb-systemd-helper >/dev/null 2>&1; then
  DEB_SYSTEMD_HELPER="deb-systemd-helper"
fi

case "$1" in
  purge)
    # Remove data and config on purge
    rm -rf /var/lib/tui-serve
    rm -rf /var/log/tui-serve
    rm -f /etc/tui-serve/env
    rm -f /etc/tui-serve/default-config.json
    rmdir /etc/tui-serve 2>/dev/null || true

    # Remove system user
    if getent passwd tui-serve >/dev/null 2>&1; then
      deluser --quiet tui-serve
    fi

    # Unmask and clean up systemd state
    if [ -n "$DEB_SYSTEMD_HELPER" ]; then
      "$DEB_SYSTEMD_HELPER" unmask tui-serve.service >/dev/null 2>&1 || true
      "$DEB_SYSTEMD_HELPER" disable tui-serve.service >/dev/null 2>&1 || true
      "$DEB_SYSTEMD_HELPER" update-state tui-serve.service >/dev/null 2>&1 || true
    else
      systemctl unmask tui-serve.service >/dev/null 2>&1 || true
      systemctl disable tui-serve.service >/dev/null 2>&1 || true
    fi
    ;;

  remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
    # Mask the unit so it can't be accidentally activated while removed
    # (but not purged — postrm purge handles that)
    if [ -n "$DEB_SYSTEMD_HELPER" ]; then
      "$DEB_SYSTEMD_HELPER" mask tui-serve.service >/dev/null 2>&1 || true
    else
      systemctl mask tui-serve.service >/dev/null 2>&1 || true
    fi
    ;;

  *)
    echo "postrm called with unknown argument \`$1'" >&2
    ;;
esac

exit 0