# packaging/debian/postrm
#!/bin/sh
# postrm: post-remove script for remote-agent-tui
# Follows dh_installsystemd canonical patterns (Debian Policy §9.3)
set -e

DEB_SYSTEMD_HELPER=""
if command -v deb-systemd-helper >/dev/null 2>&1; then
  DEB_SYSTEMD_HELPER="deb-systemd-helper"
fi

case "$1" in
  purge)
    # Remove data and config on purge
    rm -rf /var/lib/remote-agent-tui
    rm -rf /var/log/remote-agent-tui
    rm -f /etc/remote-agent-tui/env
    rm -f /etc/remote-agent-tui/default-config.json
    rmdir /etc/remote-agent-tui 2>/dev/null || true

    # Remove system user
    if getent passwd remote-agent-tui >/dev/null 2>&1; then
      deluser --quiet remote-agent-tui
    fi

    # Unmask and clean up systemd state
    if [ -n "$DEB_SYSTEMD_HELPER" ]; then
      "$DEB_SYSTEMD_HELPER" unmask remote-agent-tui.service >/dev/null 2>&1 || true
      "$DEB_SYSTEMD_HELPER" disable remote-agent-tui.service >/dev/null 2>&1 || true
      "$DEB_SYSTEMD_HELPER" update-state remote-agent-tui.service >/dev/null 2>&1 || true
    else
      systemctl unmask remote-agent-tui.service >/dev/null 2>&1 || true
      systemctl disable remote-agent-tui.service >/dev/null 2>&1 || true
    fi
    ;;

  remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
    # Mask the unit so it can't be accidentally activated while removed
    # (but not purged — postrm purge handles that)
    if [ -n "$DEB_SYSTEMD_HELPER" ]; then
      "$DEB_SYSTEMD_HELPER" mask remote-agent-tui.service >/dev/null 2>&1 || true
    else
      systemctl mask remote-agent-tui.service >/dev/null 2>&1 || true
    fi
    ;;

  *)
    echo "postrm called with unknown argument \`$1'" >&2
    ;;
esac

exit 0