# packaging/debian/postinst
#!/bin/sh
# postinst: post-install script for tui-serve
# Follows dh_installsystemd canonical patterns (Debian Policy §9.3)
set -e

# deb-systemd-invoke and deb-systemd-helper may not exist on minimal systems
# (provided by init-system-helpers, typically pre-installed)
if command -v deb-systemd-invoke >/dev/null 2>&1; then
  DEB_SYSTEMD_INVOKE="deb-systemd-invoke"
else
  DEB_SYSTEMD_INVOKE="systemctl"
fi
if command -v deb-systemd-helper >/dev/null 2>&1; then
  DEB_SYSTEMD_HELPER="deb-systemd-helper"
else
  DEB_SYSTEMD_HELPER=""
fi

json_escape() {
  printf '%s' "$1" | sed 's/\\/\\\\/g; s/"/\\"/g'
}

append_root() {
  root="$1"
  [ -n "$root" ] || return 0

  case " $ALLOWED_ROOTS " in
    *" $root "*) return 0 ;;
  esac

  ALLOWED_ROOTS="$ALLOWED_ROOTS $root"
}

generate_token() {
  if command -v openssl >/dev/null 2>&1; then
    openssl rand -base64 32 | tr '+/' '-_' | tr -d '='
  elif [ -x /usr/lib/tui-serve/node/bin/node ]; then
    /usr/lib/tui-serve/node/bin/node -e "console.log(require('crypto').randomBytes(32).toString('base64url'))"
  else
    date +%s%N | sha256sum | awk '{print $1}'
  fi
}

get_installer_home() {
  if [ -n "${SUDO_USER:-}" ] && [ "${SUDO_USER:-}" != "root" ]; then
    getent passwd "$SUDO_USER" 2>/dev/null | cut -d: -f6 || true
  fi
}

write_generated_config() {
  ALLOWED_ROOTS=""

  installer_home="$(get_installer_home)"
  if [ -n "$installer_home" ] && [ -d "$installer_home" ]; then
    append_root "$installer_home"
  fi

  append_root "/var/lib/tui-serve"

  roots_json=""
  for root in $ALLOWED_ROOTS; do
    escaped_root="$(json_escape "$root")"
    if [ -z "$roots_json" ]; then
      roots_json="\"$escaped_root\""
    else
      roots_json="$roots_json, \"$escaped_root\""
    fi
  done

  cat > /etc/tui-serve/default-config.json << EOF
{
  "commands": [
    { "id": "pi", "label": "Pi Coding Agent", "command": "pi", "allowedCwdRoots": [$roots_json] },
    { "id": "claude", "label": "Claude Coding Agent", "command": "claude", "allowedCwdRoots": [$roots_json] },
    { "id": "opencode", "label": "OpenCode", "command": "opencode", "allowedCwdRoots": [$roots_json] },
    { "id": "codex", "label": "Codex Coding Agent", "command": "codex", "allowedCwdRoots": [$roots_json] }
  ],
  "hosts": [
    { "id": "local", "name": "This Machine", "address": "localhost", "port": 5555 }
  ]
}
EOF
}

ensure_installer_home_in_config() {
  installer_home="$(get_installer_home)"
  [ -n "$installer_home" ] || return 0
  [ -d "$installer_home" ] || return 0
  [ -f /etc/tui-serve/default-config.json ] || return 0
  [ -x /usr/lib/tui-serve/node/bin/node ] || return 0

  INSTALLER_HOME="$installer_home" /usr/lib/tui-serve/node/bin/node --input-type=module <<'NODEEOF'
import { readFileSync, writeFileSync } from 'node:fs';

const configPath = '/etc/tui-serve/default-config.json';
const installerHome = process.env.INSTALLER_HOME;
const config = JSON.parse(readFileSync(configPath, 'utf8'));

// Support both array format [{id, command, ...}] and legacy object format {commandId: {bin, ...}}
const commands = Array.isArray(config.commands)
  ? config.commands
  : Object.values(config.commands ?? {});
let changed = false;

for (const command of commands) {
  if (!command || typeof command !== 'object') continue;
  const roots = Array.isArray(command.allowedCwdRoots) ? command.allowedCwdRoots : [];
  if (!roots.includes(installerHome)) {
    command.allowedCwdRoots = [installerHome, ...roots];
    changed = true;
  } else if (roots[0] !== installerHome) {
    command.allowedCwdRoots = [installerHome, ...roots.filter((root) => root !== installerHome)];
    changed = true;
  }
}

if (changed) {
  writeFileSync(configPath, JSON.stringify(config, null, 2) + '\n');
}
NODEEOF
}

auto_refresh_installer_user_service() {
  # dpkg maintainer scripts run as root, but desktop/dev installs usually use
  # the per-user service so agents inherit the real user's HOME, PATH, SSH keys,
  # tmux socket, and dotfiles. When upgrading via `sudo dpkg -i`, refresh that
  # install too so ~/.local/share/tui-serve gets the new server/web assets.
  [ -n "${SUDO_USER:-}" ] || return 0
  [ "${SUDO_USER:-}" != "root" ] || return 0
  [ -x /usr/share/doc/tui-serve/install-user-service.sh ] || return 0

  installer_user="$SUDO_USER"
  installer_home="$(getent passwd "$installer_user" 2>/dev/null | cut -d: -f6 || true)"
  installer_uid="$(id -u "$installer_user" 2>/dev/null || true)"
  [ -n "$installer_home" ] && [ -d "$installer_home" ] || return 0
  [ -n "$installer_uid" ] || return 0

  user_service="$installer_home/.config/systemd/user/tui-serve.service"
  user_prefix="$installer_home/.local/share/tui-serve"

  # Avoid surprising fresh installs: auto-refresh only when a user install or
  # copied app tree already exists. Fresh installs still print the explicit command.
  if [ ! -f "$user_service" ] && [ ! -d "$user_prefix" ]; then
    return 0
  fi

  echo ""
  echo "Refreshing TUI Serve user install for $installer_user..."

  if [ -d "/run/user/$installer_uid" ]; then
    if runuser -u "$installer_user" -- env \
      HOME="$installer_home" \
      USER="$installer_user" \
      LOGNAME="$installer_user" \
      XDG_RUNTIME_DIR="/run/user/$installer_uid" \
      /usr/share/doc/tui-serve/install-user-service.sh; then
      echo "User service assets refreshed for $installer_user."
    else
      echo "Warning: could not refresh user service automatically." >&2
      echo "Run as $installer_user:" >&2
      echo "  /usr/share/doc/tui-serve/install-user-service.sh" >&2
      echo "  systemctl --user restart tui-serve" >&2
    fi
  else
    echo "User systemd runtime not available for $installer_user; skipping auto-refresh." >&2
    echo "Run after login:" >&2
    echo "  /usr/share/doc/tui-serve/install-user-service.sh" >&2
    echo "  systemctl --user restart tui-serve" >&2
  fi
}

case "$1" in
  configure)
    # ── Create system user ──
    if ! getent passwd tui-serve >/dev/null 2>&1; then
      adduser --quiet --system --no-create-home \
        --home /var/lib/tui-serve \
        --gecos "TUI Serve" \
        --group tui-serve
    fi

    # ── Ensure data and log directories exist with correct ownership ──
    mkdir -p /var/lib/tui-serve
    mkdir -p /var/log/tui-serve
    chown -R tui-serve:tui-serve /var/lib/tui-serve
    chown -R tui-serve:tui-serve /var/log/tui-serve

    # ── Config directory ──
    mkdir -p /etc/tui-serve
    if [ -z "$2" ]; then
      write_generated_config
    elif [ ! -f /etc/tui-serve/default-config.json ]; then
      cp /usr/share/tui-serve/default-config.json \
        /etc/tui-serve/default-config.json
    fi
    ensure_installer_home_in_config
    chmod 640 /etc/tui-serve/default-config.json
    chown tui-serve:tui-serve /etc/tui-serve/default-config.json

    # ── Create env file with secure network defaults if it doesn't exist ──
    # Default bind is 0.0.0.0, so a strong generated AUTH_TOKEN is required.
    if [ ! -f /etc/tui-serve/env ]; then
      generated_token="$(generate_token)"
      cat > /etc/tui-serve/env << ENVEOF
# TUI Serve environment configuration
# Edit this file to change settings, then restart:
#   sudo systemctl restart tui-serve

# Default network bind requires authentication. For local-only no-auth mode,
# set BIND_HOST=127.0.0.1 and AUTH_TOKEN=, then restart.
AUTH_TOKEN=$generated_token
BIND_HOST=0.0.0.0
PORT=5555
NODE_ENV=production

# Paths (pre-configured for packaged installation)
TUI_SERVE_CONFIG=/etc/tui-serve/default-config.json
TUI_SERVE_DATA_DIR=/var/lib/tui-serve
TUI_SERVE_WEB_DIR=/usr/share/tui-serve/web
ENVEOF
      chmod 600 /etc/tui-serve/env
      chown tui-serve:tui-serve /etc/tui-serve/env
    else
      if ! grep -q '^BIND_HOST=' /etc/tui-serve/env; then
        echo 'BIND_HOST=0.0.0.0' >> /etc/tui-serve/env
      fi
      if grep -q '^AUTH_TOKEN=$' /etc/tui-serve/env; then
        generated_token="$(generate_token)"
        tmp_env="$(mktemp)"
        sed "s/^AUTH_TOKEN=$/AUTH_TOKEN=$generated_token/" /etc/tui-serve/env > "$tmp_env"
        cat "$tmp_env" > /etc/tui-serve/env
        rm -f "$tmp_env"
        chmod 600 /etc/tui-serve/env
        chown tui-serve:tui-serve /etc/tui-serve/env
      fi
    fi

    # ── Install banner (always shown on fresh install) ──
    if [ -z "$2" ]; then
      echo ""
      echo "============================================"
      echo " TUI Serve - Installed"
      echo "============================================"
      echo ""
      echo " Network binding uses a generated auth token by default."
      echo " System service is installed but not started by default."
      echo ""
      echo " Recommended for desktop/dev machines: install your own user service:"
      echo "   /usr/share/doc/tui-serve/install-user-service.sh"
      echo ""
      echo " Then open: http://localhost:5555"
      echo ""
      echo " For appliance/headless mode, enable the system service manually:"
      echo "   sudo systemctl enable --now tui-serve.service"
      echo ""
      echo " Auth/config:"
      echo "   sudo nano /etc/tui-serve/env"
      echo "   # Network mode: BIND_HOST=0.0.0.0 requires AUTH_TOKEN"
      echo "   # Local-only mode: BIND_HOST=127.0.0.1 may use AUTH_TOKEN="
      echo "   sudo systemctl restart tui-serve"
      echo ""
      echo " Workspace roots include: /var/lib/tui-serve"
      if [ -n "${SUDO_USER:-}" ] && [ "${SUDO_USER:-}" != "root" ]; then
        echo " Installing user's home is included when detected."
      fi
      echo " Add project roots in:"
      echo "   sudo nano /etc/tui-serve/default-config.json"
      echo ""
      echo " Config:  sudo nano /etc/tui-serve/default-config.json"
      echo " Logs:    journalctl -u tui-serve -f"
      echo " Doctor:  sudo /usr/share/doc/tui-serve/doctor.sh"
      echo "============================================"
      echo ""
    fi

    # ── systemd unit management ──
    # Unmask the unit (may have been masked by admin or previous postrm)
    if [ -n "$DEB_SYSTEMD_HELPER" ]; then
      "$DEB_SYSTEMD_HELPER" unmask tui-serve.service >/dev/null 2>&1 || true
    fi

    # ── systemd unit registration ──
    # Install the system unit but do not enable/start it by default. Desktop/dev
    # users should run the per-user service so agents inherit their real HOME,
    # PATH, tmux socket, dotfiles, and credentials. Appliance/headless installs
    # can explicitly enable the system service.
    if [ -z "${DPKG_ROOT:-}" ] && [ -d /run/systemd/system ]; then
      systemctl --system daemon-reload >/dev/null 2>&1 || true
      systemctl reset-failed tui-serve.service >/dev/null 2>&1 || true
    fi

    if [ -n "$DEB_SYSTEMD_HELPER" ]; then
      "$DEB_SYSTEMD_HELPER" update-state tui-serve.service >/dev/null 2>&1 || true
    fi

    auto_refresh_installer_user_service

    ;;

  abort-upgrade|abort-deconfigure|abort-remove)
    # Nothing to do on abort
    ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    ;;
esac

exit 0