# packaging/debian/systemd/remote-agent-tui.service
[Unit]
Description=Remote Agent TUI Manager
After=network.target tailscaled.service
Wants=tailscaled.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=remote-agent-tui
Group=remote-agent-tui
WorkingDirectory=/usr/lib/remote-agent-tui/server

# Uses bundled Node.js — no external Node.js installation required
# ExecStartPre: fast-fail on broken node binary (wrong arch, missing libs)
ExecStartPre=/usr/lib/remote-agent-tui/node/bin/node --version
ExecStart=/usr/lib/remote-agent-tui/node/bin/node dist/index.js

Restart=on-failure
RestartSec=5

# Environment — loaded from /etc/remote-agent-tui/env
# Override in /etc/systemd/system/remote-agent-tui.service.d/override.conf
EnvironmentFile=-/etc/remote-agent-tui/env
Environment=NODE_ENV=production

# Paths for packaged installation
Environment=REMOTE_AGENT_TUI_CONFIG=/etc/remote-agent-tui/default-config.json
Environment=REMOTE_AGENT_TUI_DATA_DIR=/var/lib/remote-agent-tui
Environment=REMOTE_AGENT_TUI_WEB_DIR=/usr/share/remote-agent-tui/web

StandardOutput=journal
StandardError=journal
SyslogIdentifier=remote-agent-tui

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/remote-agent-tui /var/log/remote-agent-tui
PrivateTmp=true

[Install]
WantedBy=multi-user.target